# Scripts/com.outlook-mcp.auth-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.outlook-mcp.auth-server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>node</string>
        <string>/Users/bj/Library/CloudStorage/GoogleDrive-johnsonboibeats@gmail.com/My Drive/Repos/Outlook-MCP/auth/server.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/bj/Library/CloudStorage/GoogleDrive-johnsonboibeats@gmail.com/My Drive/Repos/Outlook-MCP</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/bj/Library/Logs/Claude/outlook-auth-server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/bj/Library/Logs/Claude/outlook-auth-server-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/bj/.nvm/versions/node/v22.16.0/bin</string>
    </dict>
</dict>
</plist>
